# recipes-core/pmsrv/files/pmsrv.service
[Unit]
Description=pmsrv service

[Service]
ExecStart=/usr/bin/pmsrv
Restart=always

[Install]
WantedBy=multi-user.target
